# debianiso/var/lib/dpkg/info/rsync.list
/.
/etc
/etc/default
/etc/default/rsync
/etc/init.d
/etc/init.d/rsync
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rsync.service
/usr
/usr/bin
/usr/bin/rrsync
/usr/bin/rsync
/usr/bin/rsync-ssl
/usr/share
/usr/share/doc
/usr/share/doc/rsync
/usr/share/doc/rsync/NEWS.Debian.gz
/usr/share/doc/rsync/NEWS.md.gz
/usr/share/doc/rsync/README.Debian
/usr/share/doc/rsync/changelog.Debian.gz
/usr/share/doc/rsync/copyright
/usr/share/doc/rsync/examples
/usr/share/doc/rsync/examples/logrotate.conf.rsync
/usr/share/doc/rsync/examples/rsyncd.conf
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/rsync
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/rrsync.1.gz
/usr/share/man/man1/rsync-ssl.1.gz
/usr/share/man/man1/rsync.1.gz
/usr/share/man/man5
/usr/share/man/man5/rsyncd.conf.5.gz
/usr/share/rsync
/usr/share/rsync/scripts
/usr/share/rsync/scripts/atomic-rsync
/usr/share/rsync/scripts/cull-options
/usr/share/rsync/scripts/cvs2includes
/usr/share/rsync/scripts/file-attr-restore
/usr/share/rsync/scripts/files-to-excludes
/usr/share/rsync/scripts/git-set-file-times
/usr/share/rsync/scripts/logfilter
/usr/share/rsync/scripts/lsh
/usr/share/rsync/scripts/mnt-excl
/usr/share/rsync/scripts/munge-symlinks
/usr/share/rsync/scripts/rsync-no-vanished
/usr/share/rsync/scripts/rsyncstats
